# ansible/dockerd.conf
[program:dockerd]
command=/usr/bin/docker.io -d
user=root
autorestart=true
